# roles/backup_manager/templates/restic_backup.timer.j2
[Unit]
Description=Restic backup timer for {{ item.name }}

[Timer]
OnCalendar={{ item.schedule | default('daily') }}
Persistent=true

[Install]
WantedBy=timers.target
