# roles/pbs_client/templates/backup.timer.j2
[Unit]
Description=Run Proxmox Backup Client Daily

[Timer]
OnCalendar={{ pbs_client_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
